# .config/systemd/user/nitrogen.service
[Unit]
Description=Nitrogen background setter

[Service]
ExecStart=/usr/bin/nitrogen --restore

[Install]
WantedBy=default.target
